# templates/uitdatabank/search_api/uitdatabank-consume-queue.service.erb
[Unit]
Description=UiTdatabank search API queue consumer
Requires=elasticsearch.service
After=syslog.target network.target remote-fs.target nss-lookup.target elasticsearch.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=<%= @basedir %>
ExecStart=/usr/bin/php bin/app.php <%= @command %>
Restart=on-failure
RestartSec=10s
Type=simple

[Install]
WantedBy=multi-user.target
