# templates/uitdatabank/entry_api/uitdatabank-amqp-listener-uitpas.service.erb
[Unit]
Description=UiTdatabank AMQP listener UiTPAS queue
After=syslog.target network.target remote-fs.target nss-lookup.target

[Service]
User=www-data
Group=www-data
PIDFile=/var/run/uitdatabank-amqp-listener-uitpas.pid
WorkingDirectory=<%= @basedir %>
ExecStart=/usr/bin/php bin/udb3.php amqp-listen-uitpas
Restart=on-failure
RestartSec=10s
Type=simple

[Install]
WantedBy=multi-user.target
