# templates/uitdatabank/entry_api/uitdatabank-event-export-worker@.service.erb
[Unit]
Description=UiTdatabank event export worker %i
After=syslog.target network.target remote-fs.target nss-lookup.target
PartOf=uitdatabank-event-export-workers.target

[Service]
User=www-data
Group=www-data
PermissionsStartOnly=true
RuntimeDirectory=uitdatabank-entry-api
RuntimeDirectoryMode=0755
PIDFile=/run/uitdatabank-entry-api/uitdatabank-event-export-worker@%i.pid
WorkingDirectory=<%= @basedir %>/vendor/chrisboulton/php-resque
Environment=APP_INCLUDE=<%= @basedir %>/worker_bootstrap.php
Environment=INTERVAL=1
Environment=QUEUE=event_export
ExecStart=/usr/bin/php resque.php
Restart=on-failure
RestartSec=10s
Type=simple

[Install]
WantedBy=multi-user.target

